# profile.d/git
GIT_WORKSPACE=
GIT_PS1=
GIT_CONFIG="~/.gitconfig"
GIT_EDITOR=
GIT_MERGETOOL=
GIT_DIFFTOOL=
