# profile.d/global_complete
COLUMN=1
GWS_COMPLETION_COMMANDS_NUMBER=100