# profile.d/global_complete
$COLUMN=1
$gws_COMPLETION_COMMANDS_NUMBER=100